# root/etc/uci-defaults/luci-theme-jj
#!/bin/sh
	uci batch <<-EOF
	        set luci.themes.MyTheme=/luci-static/jj
	        commit luci
	EOF
exit 0